# Linked_Development/overlay/etc/network/if-pre-up.d/iptables
#!/bin/bash
/sbin/iptables-restore < /etc/network/iptables.up.rules
